# app/conf/etc/supervisor/conf.d/10-worker.conf
[program:worker]
command = su satis -c "/usr/bin/node /home/satis/webhook_and_worker/build-worker.js"
autostart = true
autorestart = true
priority = 10
