# app/conf/etc/supervisor/conf.d/20-webhook.conf
[program:webhook]
command = su - satis -c "/usr/bin/node /home/satis/webhook_and_worker/webhook.js"
autostart = true
autorestart = true
priority = 20
